# .config/systemd/user/clash.service
[Unit]
Description=run clash

[Service]
Type=Simple
ExecStart=clash

[Install]
WantedBy=default.target